# packaging/arch/gamerat-daemon.service
[Unit]
Description=gamerat focus-aware mouse profile daemon
Documentation=https://github.com/appulsauce/gamerat
After=graphical-session.target dbus.socket
PartOf=graphical-session.target

[Service]
Type=dbus
BusName=org.appulsauce.GameRat1
ExecStart=/usr/bin/gamerat-daemon
Restart=on-failure
RestartSec=5
# Daemon owns its own state under $XDG_CONFIG_HOME and
# $XDG_STATE_HOME; no special privileges required.

[Install]
WantedBy=graphical-session.target
