# packaging/systemd/gamerat-daemon.service
[Unit]
Description=gamerat user daemon (per-application mouse profile dispatcher)
Documentation=https://github.com/appulsauce/gamerat
After=graphical-session.target dbus.socket
PartOf=graphical-session.target

[Service]
Type=dbus
BusName=org.appulsauce.GameRat1
ExecStart=/usr/bin/gamerat-daemon
Restart=on-failure
RestartSec=3

# Hardening — daemon never reaches the network and only writes under
# $XDG_CONFIG_HOME and $XDG_STATE_HOME. ProtectHome is left at the
# default (read/write) because the daemon's own state files live
# under $HOME; ReadWritePaths narrows the surface back down.
PrivateNetwork=true
ProtectSystem=strict
ReadWritePaths=%h/.config/gamerat %h/.local/state/gamerat
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=graphical-session.target
